# launchd/in.krud.ai.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>in.krud.ai</string>
    <key>ProgramArguments</key>
    <array>
      <string>__KRUDD_BINARY__</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__KRUD_LOG_DIR__/krudd.log</string>
    <key>StandardErrorPath</key>
    <string>__KRUD_LOG_DIR__/krudd.err.log</string>
  </dict>
</plist>
